# deploy/rp-discovery.service
[Unit]
Description=Research Pipeline - Discovery Service
After=network.target
PartOf=rp-pipeline.target

[Service]
Type=simple
WorkingDirectory=/media/sam/1TB/research-pipeline
Environment=RP_DB_PATH=data/research.db
Environment=RP_DISCOVERY_PORT=8770
ExecStart=/media/sam/1TB/research-pipeline/.venv/bin/python -m services.discovery.main
Restart=always
RestartSec=5
WatchdogSec=60

PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/media/sam/1TB/research-pipeline/data

[Install]
WantedBy=rp-pipeline.target
