# deploy/rp-mcp.service
[Unit]
Description=PePeRS - MCP Server (SSE)
After=network.target rp-orchestrator.service
PartOf=rp-pipeline.target

[Service]
Type=simple
WorkingDirectory=/media/sam/1TB/pepers
Environment=RP_MCP_PORT=8776
Environment=RP_ORCHESTRATOR_URL=http://localhost:8775
Environment=RP_MCP_FLAVOR=arcade
ExecStart=/media/sam/1TB/pepers/.venv/bin/python -m services.mcp
Restart=always
RestartSec=5

PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/media/sam/1TB/pepers/data

[Install]
WantedBy=rp-pipeline.target
