# deploy/rp-orchestrator.service
[Unit]
Description=PePeRS - Orchestrator Service
# Adjust /opt/pepers paths below to match your installation directory
After=network.target
PartOf=rp-pipeline.target

[Service]
Type=simple
WorkingDirectory=/opt/pepers
Environment=RP_DB_PATH=data/research.db
Environment=RP_ORCHESTRATOR_PORT=8775
ExecStart=/opt/pepers/.venv/bin/python -m services.orchestrator.main
Restart=always
RestartSec=5
WatchdogSec=60

PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/pepers/data

[Install]
WantedBy=rp-pipeline.target
