# deploy/rp-validator.service
[Unit]
Description=PePeRS - Validator Service
After=network.target
PartOf=rp-pipeline.target

[Service]
Type=simple
WorkingDirectory=/media/sam/1TB/pepers
Environment=RP_DB_PATH=data/research.db
Environment=RP_VALIDATOR_PORT=8773
ExecStart=/media/sam/1TB/pepers/.venv/bin/python -m services.validator.main
Restart=always
RestartSec=5
WatchdogSec=60

PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/media/sam/1TB/pepers/data

[Install]
WantedBy=rp-pipeline.target
